# supervisord/ecovacs.ini
[program:ecovacs]
directory = /opt/ecovacs/
command = /opt/ecovacs/ecovacs-mqtt-gateway.py
